# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=redis-server

[program:springapp]
command=java -jar /app/app.jar

